# systemd/ptcpdump@.service
[Unit]
Description=ptcpdump packet capture on %I
Documentation=https://github.com/tobixen/bw-meter
After=network.target
# Start after the named interface appears
After=sys-subsystem-net-devices-%i.device
BindsTo=sys-subsystem-net-devices-%i.device

[Service]
Type=simple
# Create capture directory if it does not yet exist
ExecStartPre=/usr/bin/install -d -m 750 /var/lib/bw-meter/%I
# Write to a timestamped file; RuntimeMaxSec kills and restarts every 15 min
# so each run gets a fresh filename.
# Note: %%Y etc. — double-% to pass a literal % through systemd specifier expansion.
# Adjust the ptcpdump path if it is not at /usr/bin/ptcpdump.
ExecStart=/bin/sh -c 'exec /usr/bin/ptcpdump -i %I -w /var/lib/bw-meter/%I/capture-$(date +%%Y%%m%%dT%%H%%M%%S).pcapng'
RuntimeMaxSec=900
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
